# root/etc/services.d/php-fpm/run
#!/usr/bin/execlineb -P

with-contenv
php-fpm
