# root/etc/services.d/nodejs/run
#!/usr/bin/execlineb -P

s6-envdir -fn -- /tmp/cont-init-env

importas -u PKGDIR PKGDIR

foreground {
  s6-echo -- "Running 'npm run dev' for '${PKGDIR}/package.json'"
}
cd "${PKGDIR}" su-exec www-data npm run dev

